# root/etc/init.d/hp-lite
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1

PROG=/usr/bin/hp-lite
WORKDIR=/tmp/hp-lite
CONFIG=/etc/config/hp-lite
LOG_DIR=/var/log/hp-lite
CRON_TAG="# hp-lite-log-clean"
DEFAULT_RETENTION=3
LOG_FILE="$LOG_DIR/hp-lite.log"

shell_quote() {
    printf "'%s'" "$(printf "%s" "$1" | sed "s/'/'\\\\''/g")"
}

start_service() {

    mkdir -p "$LOG_DIR"
    mkdir -p "$WORKDIR"
    touch "$LOG_FILE"
    chmod 644 "$LOG_FILE"

    echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lite] Starting service..." >> "$LOG_FILE"
  
    local code=$(uci get hp-lite.global.connect_code 2>/dev/null)
    [ -z "$code" ] && {
        echo "No connect code configured, please set it in LuCI"
        echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lite] Start failed: connect_code is not configured, please set it in LuCI" >> "$LOG_FILE"
        return 1
    }
    
    local cmd="$(shell_quote "$PROG") $(shell_quote "-c=$code") 2>&1 | tee -a $(shell_quote "$LOG_FILE")"

    procd_open_instance
    procd_set_param command sh -c "$cmd"
    procd_set_param respawn
    procd_set_param cwd "$WORKDIR"
    procd_set_param file "$LOG_FILE"
    procd_close_instance 
    
    echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lite] Service started successfully, process is supervised by procd" >> "$LOG_FILE"

    local RET=$(uci get hp-lite.global.log_retention_days 2>/dev/null)
    [ -z "$RET" ] && RET=$DEFAULT_RETENTION
    case "$RET" in
        *[!0-9]*)
            RET=$DEFAULT_RETENTION
            ;;
    esac
    [ "$RET" -lt 1 ] && RET=$DEFAULT_RETENTION
    [ "$RET" -gt 3650 ] && RET=$DEFAULT_RETENTION

    local CLEAN_CMD="find '$LOG_DIR' -type f -name '*.log' -exec sh -c ': > \"\$1\"' sh {} \;"
    
    [ -f /etc/crontabs/root ] || touch /etc/crontabs/root
    sed -i "/$CRON_TAG/d" /etc/crontabs/root
    echo "0 3 */$RET * * $CLEAN_CMD $CRON_TAG" >> /etc/crontabs/root  # Run cleanup every $RET days
    /etc/init.d/cron reload >/dev/null 2>&1 || /etc/init.d/cron restart >/dev/null 2>&1
    
    echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lite] Log cleanup scheduled: clear logs every $RET days at 03:00" >> "$LOG_FILE"
}
stop_service() {
    procd_stop_instance "$PROG" 2>/dev/null
    local pids=$(pgrep -f "$PROG")
    if [ -n "$pids" ]; then
        for pid in $pids; do
            if kill -9 "$pid" 2>/dev/null; then
                echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lite] Process (PID: $pid) terminated successfully" >> "$LOG_FILE"
            else
                echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lite] Failed to terminate process (PID: $pid)" >> "$LOG_FILE"
            fi
        done
        
        pkill -9 -f "$PROG" 2>/dev/null
        echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lite] Cleaned up all hp-lite related processes" >> "$LOG_FILE"
    else
        echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lite] No running hp-lite process found" >> "$LOG_FILE"
    fi
}
